# FreeBSD/irc/qwebirc/files/qwebirc.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: qwebirc
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable qwebirc:
# qwebirc_enable (bool):	Set to "NO" by default.
#				Set it to "YES" to enable qwebirc
#
# qwebirc_ip (string):		Set to "127.0.0.1" by default.
# qwebirc_port (string):	Set to "9090" by default.

. /etc/rc.subr

name=qwebirc
rcvar=qwebirc_enable
load_rc_config $name

: ${qwebirc_enable:=NO}
: ${qwebirc_user:=%%USERS%%}
: ${qwebirc_group:=%%GROUPS%%}
: ${qwebirc_ip:=127.0.0.1}
: ${qwebirc_port:=9090}

command_interpreter=%%PYTHON_CMD%%
command=./run.py
command_args="--ip=${qwebirc_ip} --port=${qwebirc_port}"
pidfile=%%RUNDIR%%/qwebirc.pid
start_precmd=qwebirc_prestart
start_postcmd=qwebirc_poststart

required_files="%%SCRIPTDIR%%/config.py"

cd "%%SCRIPTDIR%%" || exit 1


qwebirc_prestart()
{
  if [ ! -d "${pidfile%/*}" ]
  then
    install -d -o '%%USERS%%' -g '%%GROUPS%%' "${pidfile%/*}"
  fi
}

qwebirc_poststart(){
  pid=$( check_process "${command}" "${command_interpreter}" )
  [ -z "${pid}" ] && exit 1
  printf '%s\n' "${pid}" > "${pidfile}"
}

run_rc_command "$1"
