# FreeBSD/irc/qwebirc/files/qwebirc.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: qwebird
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable qwebirc:
# qwebirc_enable (bool):	Set to "NO" by default.
#				Set it to "YES" to enable qwebirc
#

. /etc/rc.subr

name=qwebirc
rcvar=qwebirc_enable
load_rc_config $name

: ${qwebirc_enable:=NO}
: ${qwebirc_user:=%%USERS%%}
: ${qwebirc_group:=%%GROUPS%%}

command=%%PREFIX%%/bin/python2.7
command_args=run.py
pidfile=%%RUNDIR%%/ircd.pid
start_precmd=qwebirc_prestart

required_files="%%SCRIPTDIR%%/config.py"

cd "%%SCRIPTDIR%%" || exit 1


qwebirc_prestart()
{
	if [ ! -d "${pidfile%/*}" ]; then
		install -d -o '%%USERS%%' -g '%%GROUPS%%' "${pidfile%/*}"
	fi
}

run_rc_command "$1"
