# FreeBSD/net/miniature-spoon/files/miniature-spoon.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: miniature_spoon
# REQUIRE: LOGIN cleanvar
# KEYWORD: shutdown

# Configuration settings in /etc/rc.conf:
#
# miniature_spoon_enable (bool):
#   Set to "NO" by default.
#   Set it to "YES" to enable miniature_spoon
#

# miniature_spoon_config_file (str):
#   Set to "%%PREFIX%%/etc/${name}.conf" by default.
#   The master configuration file.

. /etc/rc.subr

name="miniature_spoon"
rcvar=miniature_spoon_enable

load_rc_config ${name}

: ${miniature_spoon_enable:="NO"}
: ${miniature_spoon_config_file:="%%PREFIX%%/etc/miniature-spoon.conf"}
: ${miniature_spoon_run_dir:="%%MINIATURE_SPOON_RUN%%"}

start_precmd="${name}_precmd"
restart_precmd="${name}_checkconfig"
stop_postcmd="${name}_poststop"

#extra_commands='initialise reload configtest'

pidfile="/var/run/miniature-spoon.pid"

required_files="${miniature_spoon_config_file}"

_miniature_spoon_program="%%PREFIX%%/sbin/miniature-spoon"
#_miniature_spoon_arguments="--quiet --config-file=${miniature_spoon_config_file}"
_miniature_spoon_arguments="${miniature_spoon_config_file}"

command="/usr/sbin/daemon"
command_args="-f -P ${pidfile} -r -- ${_miniature_spoon_program} ${_miniature_spoon_arguments}"
#procname=${_miniature_spoon_program} # do not override procname as program runs under daemon


miniature_spoon_checkconfig()
{
  #echo "Performing sanity check on miniature_spoon configuration:"
  #eval ${command} ${miniature_spoon_flags} -t
}

miniature_spoon_precmd()
{
  miniature_spoon_checkconfig
}

miniature_spoon_poststop()
{
    rm -f "${pidfile}"
}

run_rc_command "$1"
